# reader-distribution-debian/src/main/deb/control/postrm
#!/bin/sh

set -e

# Remove symlink.
rm -f /usr/bin/reader

# Remove startup scripts.
update-rc.d -f reader remove
