# reader-distribution-debian/src/main/deb/control/prerm
#!/bin/sh

set -e

# Stop Reader service.
if [ -e /etc/init.d/reader ]; then
  invoke-rc.d reader stop
fi
